# media-tv/kodi/files/kodi-initd
#!/sbin/openrc-run
# Distributed under the terms of the GNU General Public License v2

description="Start and stop Kodi"

EXEC="/usr/bin/kodi-standalone"

depend() {
	need net
	need localmount
	need local
}

start() {
	ebegin "Starting kodi"
	if [ -z $KODI_USER ]; then
			eend "No user defined in /etc/conf.d/kodi"
	fi

	start-stop-daemon --start --quiet --chuid $KODI_USER --background --exec $EXEC -- $KODI_OPTS

	eend $?
}

stop() {
	ebegin "Stopping kodi"

	for p in $(ps ax | grep kodi | grep standalone | cut -d" " -f1); do
		kill -9 $p
	done

 	eend 0
}

